# training-core/src/main/resources/profile.d/hadoop.sh
HADOOP_HOME=/usr/local/service/hadoop
export HADOOP_CONF_DIR=${HADOOP_HOME}/etc/hadoop
export HADOOP_CLASSPATH=`${HADOOP_HOME}/bin/hadoop classpath`