# training-core/src/main/resources/profile.d/mvn.sh
MAVEN_HOME=/usr/local/service/maven
PATH=${PATH}:${MAVEN_HOME}/bin

export MAVEN_HOME
export PATH